# packaging/launchd/com.discstream.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.discstream.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>mkdir -p "$HOME/Library/Application Support/DiscStream/logs"; set -a; . "$HOME/Library/Application Support/DiscStream/config/discstream.env"; set +a; cd "/path/to/DiscStream"; exec pnpm --filter @discstream/server start >> "$HOME/Library/Application Support/DiscStream/logs/server.out.log" 2>> "$HOME/Library/Application Support/DiscStream/logs/server.err.log"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/discstream.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/discstream.launchd.err.log</string>
</dict>
</plist>
